# led-controller.service
[Unit]
Description=LED Controller Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/led_controller
ExecStart=/usr/bin/python3 /home/pi/led_controller/led_controller.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
